# pkgs/gitea-actions-vm-image/runner-cycle.service
[Unit]
Description=Cycle the Gitea Actions runner (SIGTERM main pid; drains in-flight job, then exits)

[Service]
Type=oneshot
# --kill-whom=main is load-bearing: the default (all) would SIGTERM every
# process in the cgroup, including the running job's children, instead of
# letting act_runner drain its own job and exit cleanly. With main only,
# act_runner's signal handler cancels polling, waits up to
# runner.shutdown_timeout for the in-flight job to finish, then exits;
# gitea-runner.service's ExecStopPost=poweroff -f halts the VM and the
# host's Restart=always brings up a fresh ephemeral registration.
ExecStart=-/bin/systemctl kill --kill-whom=main --signal=SIGTERM gitea-runner.service
